# docker/supervisord/conf.d/programm.conf
[program:consumer]
process_name=%(program_name)s_%(process_num)02d
command=/root/consumer
autostart=true
autorestart=true
numprocs=5
stdout_logfile=/root/consumer.log